# ops/plumbline-faucet.service
[Unit]
Description=Plumbline Calibration Faucet (tFIL + USDFC)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=faucet
Group=faucet
WorkingDirectory=/opt/plumbline-faucet
EnvironmentFile=/opt/plumbline-faucet/.env
ExecStart=/usr/bin/env pnpm start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/opt/plumbline-faucet
# Faucet only ever opens outbound HTTP for RPC + Turnstile and listens
# on its configured port; restrict accordingly.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
